# config/autostart/Emacs.desktop
[Desktop Entry]
Name=emacs
Icon=emacs
Exec=systemctl --user restart emacs
Terminal=false
Type=Application
