# config/autostart/Ulauncher.desktop
[Desktop Entry]
Name=ulauncher
Icon=ulauncher
Exec=ulauncher --hide-window
Terminal=false
Type=Application
